# deploy/com.googleplay-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.googleplay-mcp</string>

    <!-- 请将 uv 路径替换为实际路径 (which uv), 项目路径和配置文件路径也需替换 -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/uv</string>
        <string>--directory</string>
        <string>/path/to/googleplay-mcp</string>
        <string>run</string>
        <string>python</string>
        <string>-m</string>
        <string>googleplay_mcp</string>
        <string>--config</string>
        <string>/path/to/googleplay-mcp/config.yaml</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/var/log/googleplay-mcp.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/googleplay-mcp.error.log</string>

    <key>WorkingDirectory</key>
    <string>/path/to/googleplay-mcp</string>
</dict>
</plist>
